# prayer-player.service
[Unit]
Description=Prayer Music Player
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ahmed
WorkingDirectory=/home/ahmed/YouTubeMusicPrayerPlayer
ExecStart=/usr/bin/python3 /home/ahmed/YouTubeMusicPrayerPlayer/prayer_player.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
